# supervisord.conf
[supervisord]

[program:ruuvitag-mqtt-bridge]
command=./ENV/bin/python ruuvitag-mqtt-bridge.py

